# config/com.tamm.wednesday-tts.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.tamm.wednesday-tts</string>

    <key>ProgramArguments</key>
    <array>
        <string>/REPLACE_WITH_VENV_PATH/bin/python</string>
        <string>-m</string>
        <string>wednesday_tts.server.daemon</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/REPLACE_WITH_REPO_PATH</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/tmp/wednesday-tts.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/wednesday-tts.err</string>
</dict>
</plist>
